# service/mihomo.service
[Unit]
Description=Mihomo (Clash Meta) Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/mihomo -f /etc/mihomo/config.yaml
Restart=always
RestartSec=3
LimitNOFILE=1048576
WorkingDirectory=/etc/mihomo

[Install]
WantedBy=multi-user.target
